# debian/DEBIAN/prerm
#!/bin/bash

set -e
 
# stop the server
invoke-rc.d meteor-desktop stop
